# docker/build/php/systemd/app.target
[Unit]
Description=App
Type=simple
PIDFile=/work/www/myunit/shared/tmp/pids/service.pid